# misc/py-diamcp/pkg-plist
etc/rc.d/diamcp
%%PYTHON_SITELIBDIR%%/diamcp/requirements.txt
%%PYTHON_SITELIBDIR%%/diamcp/server.py
%%PYTHON_SITELIBDIR%%/diamcp/tools/__init__.py
%%PYTHON_SITELIBDIR%%/diamcp/tools/base.py
%%PYTHON_SITELIBDIR%%/diamcp/tools/brave.py
%%PYTHON_SITELIBDIR%%/diamcp/tools/dux.py
%%PYTHON_SITELIBDIR%%/diamcp/tools/builtin.py
